# systemd/linux-agent-telegram.service
[Unit]
Description=Linux Mini Agent — Telegram Bot
After=network.target linux-agent-listen.service
Wants=network-online.target linux-agent-listen.service

[Service]
Type=simple
User=__USER__
Group=__USER__
WorkingDirectory=__REPO_ROOT__/apps/telegram
ExecStart=__UV_PATH__ run python main.py
Restart=on-failure
RestartSec=5
EnvironmentFile=__REPO_ROOT__/.env
Environment=VIRTUAL_ENV=
Environment=DISPLAY=__DISPLAY__
Environment=HOME=__HOME_DIR__
Environment=PATH=__HOME_DIR__/.local/bin:__HOME_DIR__/.npm-global/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=linux-agent-telegram

[Install]
WantedBy=multi-user.target
