# systemd/linux-agent-listen.service
[Unit]
Description=Linux Mini Agent — Listen Server (port 7600)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
Group=__USER__
WorkingDirectory=__REPO_ROOT__/apps/listen
ExecStart=__UV_PATH__ run python main.py
Restart=on-failure
RestartSec=5
EnvironmentFile=__REPO_ROOT__/.env
Environment=VIRTUAL_ENV=
Environment=DISPLAY=:0
Environment=HOME=__HOME_DIR__
Environment=PATH=__HOME_DIR__/.local/bin:__HOME_DIR__/.npm-global/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=linux-agent-listen

[Install]
WantedBy=multi-user.target
